# homesite/scripts/homesite.service
[Unit]
Description=Homesite main service
After=network-online.target

[Service]
Type=forking
EnvironmentFile=/etc/environment
ExecStart=/usr/bin/python3 /opt/homesite/homesite/scripts/control.py start
ExecStop=/usr/bin/python3 /opt/homesite/homesite/scripts/control.py stop
ExecReload=/usr/bin/python3 /opt/homesite/homesite/scripts/control.py reload
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
